# archiso/airootfs/opt/exam-app/examapp.service
[Unit]
Description=Secure Examination Application
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=%I
Environment=DISPLAY=:0
ExecStart=/usr/bin/python3 /opt/examapp/main.py --config /etc/examapp/config.json
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
